# expressvpn/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/default/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/evpn
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/expressvpn
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/expressvpn/
-rwxr-xr-x	root/root	usr/lib/expressvpn/expressvpn-agent.desktop
-rwxr-xr-x	root/root	usr/lib/expressvpn/icon.png
-rwxr-xr-x	root/root	usr/lib/expressvpn/libxvclient.so
-rwxr-xr-x	root/root	usr/lib/expressvpn/openvpn
-rwxr-xr-x	root/root	usr/lib/expressvpn/version
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/expressvpnd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rwxr-xr-x	root/root	usr/share/man/man1/expressvpn.1.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/expressvpn/
drwxr-xr-x	root/root	var/lib/expressvpn/certs/
